# src/oomd/etc/oomd.service.in
[Unit]
Description=Userland out-of-memory killer daemon
Documentation=man:oomd(1) https://github.com/facebookincubator/oomd#configuration
After=system.slice
ConditionFileIsExecutable=@bindir@/oomd

[Service]
Environment="DAEMON_OPTS=--interval 1 --config @oomdconfdir@/oomd.json"
EnvironmentFile=-@envfiledir@/oomd
ExecStart=@bindir@/oomd $DAEMON_OPTS
Restart=always
SyslogIdentifier=oomd
MemoryLow=64M

[Install]
WantedBy=multi-user.target
