# src/oomd/etc/oomd.service.in
[Unit]
Description=Userland out-of-memory killer daemon
After=system.slice

[Service]
ExecStart=@bindir@/oomd --interval 1 --config @oomdconfdir@/oomd.json
Restart=always
SyslogIdentifier=oomd
MemoryLow=64M

[Install]
WantedBy=multi-user.target
